# deploy/systemd/skyguardian-alert-relay.service
[Unit]
Description=SkyGuardian Telegram Alert Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SkyGuardianUa
ExecStart=/usr/bin/php artisan alerts:relay --sleep=1
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGTERM
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/SkyGuardianUa/storage /var/www/SkyGuardianUa/bootstrap/cache
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skyguardian-alert-relay

[Install]
WantedBy=multi-user.target
