# deploy/systemd/skyguardian-alerts-v1.service
[Unit]
Description=SkyGuardian v1 alerts worker
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/skyguardian
ExecStart=/usr/bin/php /var/www/skyguardian/bin/alerts-worker.php
Restart=always
RestartSec=5
User=www-data
Group=www-data
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
