# deploy/systemd/skyguardian-news-relay.service
[Unit]
Description=SkyGuardian News Telegram Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SkyGuardianUa
ExecStart=/usr/bin/php /var/www/SkyGuardianUa/artisan news:relay
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/SkyGuardianUa/storage /var/www/SkyGuardianUa/bootstrap/cache
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skyguardian-news-relay

[Install]
WantedBy=multi-user.target
