# deploy/skyguardian-data-channels.service
[Unit]
Description=SkyGuardian Telegram data channel monitor
After=network-online.target php8.3-fpm.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SkyGuardianUa/storage/telegram-sessions
ExecStart=/usr/bin/php /var/www/SkyGuardianUa/bin/data-channel-worker.php
Restart=always
RestartSec=5
RuntimeMaxSec=600
TimeoutStopSec=20
KillMode=control-group
KillSignal=SIGTERM
SendSIGKILL=yes
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
