# deploy/skyguardian-data-news.service
[Unit]
Description=SkyGuardian Telegram news data worker
After=network-online.target php8.3-fpm.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SkyGuardianUa
ExecStartPre=/usr/bin/mkdir -p /var/www/SkyGuardianUa/storage/telegram-news-sessions
ExecStart=/bin/sh -c 'while true; do /usr/bin/php /var/www/SkyGuardianUa/bin/data-channel-worker.php news; sleep 1; done'
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
KillMode=control-group
SendSIGKILL=yes
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
